# conf.d/fluxbox.conf
[program:fluxbox]
command=fluxbox
environment=DISPLAY=:1
autorestart=true
