# conf.d/xterm.conf
[program:xterm]
command=xterm
environment=DISPLAY=:1
autorestart=true
